# files/etc/zksr/supervisor.conf
; Supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[supervisord]
nodaemon=true
directory=/tmp
pidfile=/tmp/zksr-supervisor.pid
logfile=/dev/null
logfile_maxbytes=0

[program:zookeeper]
command=%(ENV_ZOOKEEPER_HOME)s/bin/zkServer.sh start-foreground
environment=JVMFLAGS="-Dlog4j.configuration=file:/etc/zksr/zookeeper-log4j.properties"
priority=1
startsecs=5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:kafka]
command=sh /usr/local/bin/start-kafka.sh
priority=2
startsecs=5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:schema-registry]
command=sh /usr/local/bin/start-schema-registry.sh
priority=3
startsecs=5
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
